# petit_furata/furuta_python/autostart/disable_hdmi.desktop
[Desktop Entry]
Type=Application
Name=disable_hdmi
NoDisplay=true
Exec=lxterminal -e /home/serge/hdmi_autostart.sh
