# furuta_python/autostart/furuta.desktop
[Desktop Entry]
Type=Application
Name=Furuta
NoDisplay=true
Exec=/home/serge/furuta_autostart.sh
