# furuta_python/autostart/pigpiod.desktop
[Desktop Entry]
Type=Application
Name=pigpio autostart
Exec=sudo pigpiod -s 5
